# leaf-dashboard/leaf-dashboard.service
[Unit]
Description=Leaf Dashboard Server
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/var/www/dashboard
ExecStart=/usr/bin/node dashboard-server.js
Restart=always
RestartSec=10
Environment=NODE_ENV=production

[Install]
WantedBy=multi-user.target 